# systemd/fstrim-weekly.service
[Unit]
Description=Weekly fstrim for mounted filesystems

[Service]
Type=oneshot
ExecStart=/sbin/fstrim -av
